# taysr.service
[Unit]
Description=Taysr Discord Bot Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/taysr
Environment=NODE_ENV=production
EnvironmentFile=/home/YOUR_USERNAME/taysr/.env
ExecStart=/usr/bin/node /home/YOUR_USERNAME/taysr/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taysr

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/YOUR_USERNAME/taysr

[Install]
WantedBy=multi-user.target
